# vendor/package/pppd-scripts/files/pppd.service
[Unit]
Description=PPP over Serial link
After=network.target

[Service]
ExecStart=/usr/bin/pon wcdma
ExecStop=/usr/bin/poff
Restart=on-failure

[Install]
WantedBy=multi-user.target
