# vendor/package/init-scripts/files/50-usb.rules
SUBSYSTEM=="extcon", ACTION=="change", ENV{DEVTYPE}=="extcon0", RUN+="/usr/bin/usb_change.sh '$env{STATE}'"